# cronjobs-in-Arch/rotate-connectivity-check-log.service
[Unit]
Description=Rotates the connectivity-checks log file

[Service]
Type=simple
ExecStart=/home/shaun/scripts/rotatelog.sh
